# target/linux/mediatek/mt7981/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
        /sbin/smp.sh
        sleep 2
        ifconfig wan down
        sleep 2
        ifconfig eth1 down
        ifconfig eth1 up
        sleep 2
        ifconfig wan up
}
